# deploy/config/cron.hourly/00-drupal-cron
#!/bin/bash

## No longer used - this site is now static

# /usr/bin/wget -O - -q -t 1 http://www.soundsoftware.ac.uk/cron.php

