# deploy/config/cron.hourly/00-drupal-cron
#!/bin/bash
/usr/bin/wget -O - -q -t 1 http://www.soundsoftware.ac.uk/cron.php

