# deploy/config/cron.hourly/20-convert-external-repos
#!/bin/bash
sudo -H -u www-data /var/www/code/reposman/run-external.sh
