# deploy/config/cron.hourly/30-expire-explore-cache
#!/bin/sh

rm -f /var/www/code/tmp/cache/*/*/views*explore*
