# deploy/config/cron.daily/10-extract-docs
#!/bin/bash

cd /tmp

/var/www/code/docgen/extract-docs.sh

