# deploy/config/cron.daily/20-check-end-of-external-repo-log
#!/bin/bash

tail -2 /var/log/external-repos.log
